# external/custom/board/rk3568/rootfs-overlay/etc/systemd/system/reset-usb2phy.service
[Unit]
Description=Reset USB2 PHY for OTG port hotplug fix
After=systemd-modules-load.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reset-usb2phy
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
