# external/custom/board/rk3568/rootfs-overlay/etc/systemd/system/update-issue.timer
[Unit]
Description=Update system banner at boot

[Timer]
OnBootSec=30s
# NOTE: Removed OnUnitActiveSec - periodic refresh was restarting getty/sway

[Install]
WantedBy=timers.target
