# external/custom/board/rk3568/rootfs-overlay/etc/systemd/system/update-issue.service
[Unit]
Description=Update /etc/issue with system info

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-issue
