# external/custom/board/rk3568/rootfs-overlay/etc/systemd/system/iwd.service.d/wait-for-modules.conf
[Unit]
# Ensure WiFi modules are loaded before iwd starts
After=systemd-modules-load.service
